# conf/network/copy-wlan.service
[Unit]
Before=network.target

[Service]
ExecStart=/usr/local/bin/copy-wlan.sh

[Install]
WantedBy=default.target